# buildfiles/etcdk
#!/sbin/openrc-run

command="/opt/etcdkeeper/etcdkeeper"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
